# itire.ru/crontab
*/3     *       *       *       *       /usr/local/bin/php /home/alex/data/www/2.mosharov.com/xml_run.php
